# .config/autostart/OpenRGB.desktop
[Desktop Entry]
Categories=Utility;
Comment=OpenRGB 0.71, for controlling RGB lighting.
Icon=OpenRGB
Name=OpenRGB
StartupNotify=true
Terminal=false
Type=Application
Exec=env QT_QPA_PLATFORM=wayland /home/shreyas/.config/autostart/openrgb.sh
